# packaging/fedora/engine-service.sysv.in
#! /bin/sh
#
# Start/stop oVirt Engine
#
# chkconfig: - 65 34
# description: oVirt Engine
# pidfile: /var/run/ovirt-engine.pid

### BEGIN INIT INFO
# Provides: ovirt-engine
# Short-Description: oVirt Engine
### END INIT INFO

# Source function library.
. /etc/rc.d/init.d/functions

[ -f /etc/sysconfig/ovirt-engine ] && . /etc/sysconfig/ovirt-engine

RETVAL=0
PROG="oVirt Engine"
LOCKFILE=/var/lock/subsys/engine-service
PIDFILE=${PIDFILE:-"/var/run/ovirt-engine.pid"}
USER="@ENGINE_USER@"

# See how we were called.
case "$1" in
  start)
	if [ `id -u` -ne 0 ]; then
		echo $"Insufficient privilege" 1>&2
		exit 4
	fi
	echo -n $"Starting $PROG: "
	ulimit -n ${FILENO:-65535}
	touch "${PIDFILE}"
        chown "${USER}" "${PIDFILE}"
	daemon --user "${USER}" --pidfile="${PIDFILE}" \
		/usr/bin/engine-service --pidfile="${PIDFILE}" start
	RETVAL=$?
	echo
	if [ $RETVAL -eq 0 ]; then
		touch $LOCKFILE
	else
		if [ -f $LOCKFILE ]; then
			RETVAL=0
		fi
	fi
	;;
  stop)
	echo -n $"Stopping $PROG: "
	killproc -p "${PIDFILE}"
	RETVAL=$?
	echo
	[ $RETVAL -eq 0 ] && rm -f $LOCKFILE
	;;
  status)
	status -p "${PIDFILE}"
	RETVAL=$?
	;;
  restart)
	$0 stop
	$0 start
	RETVAL=$?
	;;
  *)
	echo $"Usage: $0 {start|stop|status|restart}"
	exit 2
esac

exit $RETVAL
